# docker/nessus/supervisor-nessus.ini
[program:nessus]
command=/etc/init.d/nessusd start
redirect_stderr=false
stdout_logfile=/var/log/nessus/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/var/log/nessus/stderr
stderr_logfile_maxbytes=0
stopsignal=INT
numprocs_start=1
startsecs=2
autostart=true
autorestart=true
